# roles/elk_indexer/templates/supervisor-cerebro.ini.j2
[program:cerebro]
command={{ cerebro_path }}/cerebro-{{ cerebro_ver }}/bin/cerebro
autostart=true
autorestart=true
stderr_logfile=/var/log/cerebro.err.log
stdout_logfile=/var/log/cerebro.out.log
